# qinglong/rootfs/etc/supervisord.conf
[program:qinglong]
command=python3 -m qinglong
autostart=true
autorestart=true
stderr_logfile=/var/log/qinglong.err.log
stdout_logfile=/var/log/qinglong.out.log
